# doc/Autres rendus/NVIDIA_database_save/Le_Corronc_Ronan/dump-database.timer
[Unit]
Description=Dump database

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
