# doc/Autres rendus/NVIDIA_database_save/Le_Corronc_Ronan/dump-database.service
[Unit]
Description=Dump database

[Service]
ExecStart=/usr/bin/dump-database.sh

[Install]
WantedBy=multi-user.target
